# systemd-services/user/check_mail.service
[Unit]
Description=Check for new mail and notify

[Service]
Type=oneshot
ExecStart=/home/peter/.local/bin/new_mail_offlineimap.sh
