# systemd-services/cool.service
[Unit]
Description=Run script to slow down fan

[Service]
ExecStart=/home/peter/bin/cool
